# examples/oci-a1-flex-hunter@.service
# Example only. Instance %i selects a protected environment and isolated state directory.
[Unit]
Description=OCI A1 Flex Hunter target %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=oci-a1-hunter
Group=oci-a1-hunter
WorkingDirectory=/opt/oci-a1-flex-hunter
EnvironmentFile=/etc/oci-a1-flex-hunter/%i.env
ExecStart=/opt/oci-a1-flex-hunter/.venv/bin/oci-a1-flex-hunter run --continuous --max-attempts 5
TimeoutStartSec=20min
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6
ReadWritePaths=/var/lib/oci-a1-flex-hunter/%i

# Deliberately no [Install] section: enable only the paired timer after review.
